# releng/serverconf/etc/systemd/system/docker.service.d/docker.conf.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd --bip={{ docker_bridge_ip }}